# docker/crontab
@hourly php -f /var/www/html/utils/fetch_reports.php
